# files/Proxmox-VE/mnt-backups\x2dremote.automount
[Unit]
Description=AutoMount for /mnt/backups-remote
[Automount]
Where=/mnt/backups-remote
TimeoutIdleSec=600
[Install]
WantedBy=multi-user.target